# infra/systemd/rpb-bot.service
[Unit]
Description=RPB Discord Bot (Bun runtime)
After=network.target network-online.target redis-server.service postgresql@18-main.service
Wants=network-online.target redis-server.service postgresql@18-main.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/rpbey/apps/bot

# Env (partage le .env du dashboard)
EnvironmentFile=/home/ubuntu/rpbey/apps/bot/.env
Environment=NODE_ENV=production

# Bun runtime tuning (retuné 2026-04-23 — 43 GB RAM dispo, 12 CPU)
# BUN_CONFIG_MAX_HTTP_REQUESTS : bot gacha/tournaments font beaucoup d'appels Discord REST
Environment=BUN_CONFIG_MAX_HTTP_REQUESTS=4096
Environment=BUN_RUNTIME_TRANSPILER_CACHE_PATH=/home/ubuntu/rpbey/apps/bot/.bun-cache
Environment=BUN_JSC_forceRAMSize=2147483648
Environment=UV_THREADPOOL_SIZE=24

# Vérification soft du bundle gacha-client (sert /play/* quand présent).
# Idempotent ; n'échoue pas si bundle absent (api-server retourne 503).
# ExecStartPre disabled — script removed post-monorepo refactor

# Bun runs compiled JS (SWC — bot was `bot/dist/index.js` in monorepo rpb-dashboard,
# now `dist/index.js` relative to apps/rpb-bot since bot/ contents moved up one level)
ExecStart=/home/ubuntu/.bun/bin/bun dist/index.js

# Graceful shutdown
TimeoutStopSec=120
KillMode=mixed
KillSignal=SIGTERM

# Restart policy
Restart=always
RestartSec=2

# Memory & resource limits — retuné 2026-04-30 post PG18 local (RSS mesuré ~298 MiB).
# Bot Discord low-memory ; pas de bundling, juste discordx + Prisma + canvas occasional.
# JSC(2G) reste pour heap JIT ; High(1G) déclenche GC bien avant Max.
# Libère ~2 GiB pour PostgreSQL local.
MemoryMax=2G
MemoryHigh=1G
# MemoryLow : 256M garantis hors-reclaim (RSS prod ~298M, cold path Discord events).
MemoryLow=256M
MemorySwapMax=512M
TasksMax=1024
# Bot REST API prod (scrapers Vercel) : protéger hors OOM kill
OOMScoreAdjust=-500
# cgroup v2 weights — boost CPU+IO vs jobs périodiques (batch). Default cgroup = 100.
CPUWeight=200
IOWeight=200

# Security hardening (JIT-safe — see infra/systemd/snippets/hardening-jit.conf)
# ⚠️ DO NOT enable MemoryDenyWriteExecute — breaks Bun/V8 JIT.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@obsolete @privileged
ReadWritePaths=/home/ubuntu/rpbey/apps/bot/.bun-cache /home/ubuntu/rpbey/apps/bot/data

# File descriptor limits
LimitNOFILE=524288

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpb-bot

[Install]
WantedBy=multi-user.target
