# infra/systemd/rpbey-embed.service
# rpbey-embed.service — sidecar d'embeddings/rerank local (monorepo rpbey)
# Process ONNX (transformers.js → onnxruntime natif) isolé du bundle Next/web.
# Bind loopback only (127.0.0.1:7077) ; le web/bot ne parle qu'en fetch HTTP local.
# Profil ressources conforme à docs/rag-unified-pattern.md §5 (yield à l'interactif).
# Endpoints : GET /health · POST /embed · POST /rerank.
[Unit]
Description=RPBEY Embed/Rerank sidecar (e5-small + bge-reranker, Bun runtime)
After=network.target network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/rpbey

EnvironmentFile=-/home/ubuntu/rpbey/apps/embed-sidecar/.env
Environment=NODE_ENV=production
Environment=EMBED_PORT=7077
# Cache modèles ONNX maîtrisé — doit figurer dans ReadWritePaths (ProtectSystem=strict).
Environment=EMBED_CACHE_DIR=/home/ubuntu/.cache/rpbey-embed-models
Environment=BUN_RUNTIME_TRANSPILER_CACHE_PATH=/home/ubuntu/rpbey/apps/embed-sidecar/.bun-cache

ExecStart=/home/ubuntu/.bun/bin/bun apps/embed-sidecar/server.ts

TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM
Restart=always
RestartSec=3

# Ressources (§5 du contrat RAG : e5 ≈ 500 Mo + bge-reranker ≈ 400 Mo ≈ 1 Go chaud,
# plafond 3 Go ; cède le CPU/IO à l'interactif — poids par défaut 100).
MemoryMax=3G
MemoryHigh=2560M
MemoryLow=256M
TasksMax=2048
OOMScoreAdjust=-100
CPUWeight=30
IOWeight=30
Nice=10

# Hardening (JIT-safe — PAS de MemoryDenyWriteExecute, casse le JIT Bun/onnxruntime)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@obsolete @privileged

# Loopback only (§5) — aucun trafic réseau hors localhost.
IPAddressDeny=any
IPAddressAllow=localhost

# Écriture : cache modèles ONNX (1er download) + cache transpileur Bun.
ReadWritePaths=/home/ubuntu/.cache/rpbey-embed-models /home/ubuntu/rpbey/apps/embed-sidecar/.bun-cache

LimitNOFILE=131072
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpbey-embed

[Install]
WantedBy=multi-user.target
