# apps/gacha-server/deploy/rpbey-gacha.service
# rpbey-gacha.service — serveur de jeu gacha (Colyseus 0.17 sur Bun), :5050.
# REST économie consommée par le bot (127.0.0.1:5050) + Discord Activity
# (token-exchange /discord_token + GachaRoom temps réel). Backé par la DB
# partagée @rpbey/db (socket local, défauts du client — aucune var DB requise).
#
# Install :
#   sudo cp apps/gacha-server/deploy/rpbey-gacha.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now rpbey-gacha
[Unit]
Description=RPBey Gacha Server (Colyseus / Bun runtime)
After=network.target network-online.target postgresql@18-main.service
Wants=network-online.target postgresql@18-main.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/rpbey/apps/gacha-server

# Partage le .env du bot (AUTH_SECRET pour le JWT Colyseus de la Room).
# La connexion DB passe par le socket local via les défauts de @rpbey/db.
EnvironmentFile=-/home/ubuntu/rpbey/apps/bot/.env
Environment=NODE_ENV=production
Environment=GACHA_HOST=127.0.0.1
Environment=GACHA_PORT=5050
Environment=BUN_RUNTIME_TRANSPILER_CACHE_PATH=/home/ubuntu/rpbey/apps/gacha-server/.bun-cache

# Bun exécute le TS directement (Colyseus schema = décorateurs ; tsconfig active
# experimentalDecorators + emitDecoratorMetadata). Pas d'étape de build.
ExecStart=/home/ubuntu/.bun/bin/bun src/index.ts

TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM
Restart=always
RestartSec=2

# Ressources — serveur léger (REST + WS), RSS attendu ~150 MiB.
MemoryMax=1G
MemoryHigh=512M
TasksMax=512
OOMScoreAdjust=-300
CPUWeight=150
IOWeight=150

# Durcissement (JIT-safe — NE PAS activer MemoryDenyWriteExecute : casse le JIT Bun).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@obsolete @privileged
ReadWritePaths=/home/ubuntu/rpbey/apps/gacha-server/.bun-cache

LimitNOFILE=262144

StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpbey-gacha

[Install]
WantedBy=multi-user.target
